# lib/crontab
*/1 * * * * bundle exec whenever --update-crontab && bundle exec rake update_counts:update >> /var/log/cron.log 2>&1
